# packages/noticeboard/noticeboard-0.0.1.tar.gz/noticeboard-0.0.1/chores.crontab
# run my daily tasks, well before the earliest I'm likely to be up,
# but on the right day, and not excessively early (to get the weather
# fairly soon before the day):
0 4 * * * $HOME/open-projects/github/hillwithsmallfields/noticeboard/chores
