# packages/prometheus-disk-exporter/prometheus-disk-exporter-0.0.4.tar.gz/prometheus-disk-exporter-0.0.4/systemd/prometheus-disk-exporter.service
[Unit]
Description=Service to export S.M.A.R.T. and partition data as Prometheus metrics
Documentation=https://github.com/pbogre/prometheus-disk-exporter

[Service]
ExecStart=/opt/prometheus-disk-exporter/bin/prometheus_disk_exporter
Restart=always
User=prometheus

[Install]
WantedBy=multi-user.target
